# deploy/digitalocean/pos-file.service
[Unit]
Description=POS File v3
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
WorkingDirectory=/opt/pos/file-v3
Environment=NODE_ENV=production
Environment=PATH=/opt/node22/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/opt/node22/bin/node dist/main.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
